# research/deploy/systemd/alpha-lab-d-pre.timer
# alpha-lab D 盤前報告 trigger timer — 美東 Mon-Fri 09:00
# Timezone=America/New_York 讓 systemd 自動處理 DST(EST 為 UTC-5,EDT 為 UTC-4)
# OnCalendar 只給 HH:MM:SS,在 Timezone 區間當地時間觸發

[Unit]
Description=alpha-lab D pre-market workflow trigger timer (Mon-Fri 09:00 ET)

[Timer]
OnCalendar=Mon..Fri 09:00:00
Timezone=America/New_York
Persistent=true
Unit=alpha-lab-d-pre.service

[Install]
WantedBy=timers.target
