# research/deploy/systemd/alpha-lab-d-pre.service
# alpha-lab D workflow 盤前報告 trigger — 美東 Mon-Fri 09:00
#
# 美東時間 DST-aware:用 Timezone=America/New_York 而不是固定 UTC,
# 讓 systemd 自動處理 EST/EDT 切換 (09:00 ET = 13:00/14:00 UTC)。

[Unit]
Description=alpha-lab trigger D pre-market workflow (US 09:00 ET weekdays)
After=network-online.target alpha-lab-workflow.service

[Service]
Type=oneshot
# User 跑時期由部署腳本 sed 改成實際 VM user(gcloud VM 預設非 alpha)
User=__PUBLISH_USER__
WorkingDirectory=/opt/alpha-lab/research
EnvironmentFile=/opt/alpha-lab/research/.env
ExecStart=/usr/bin/curl -sS -X POST -H 'Content-Type: application/json' --max-time 30 http://127.0.0.1:8090/d/pre

[Install]
WantedBy=timers.target
