# research/deploy/systemd/alpha-lab-a.timer
# alpha-lab A workflow trigger timer — 每 6 小時一次
# OnCalendar=`*-*-* 00/6:00:00` 對應 00:00 / 06:00 / 12:00 / 18:00

[Unit]
Description=alpha-lab A workflow trigger timer (every 6h)

[Timer]
OnCalendar=*-*-* 00/6:00:00
# 如果錯過(VM 關機期間),補跑一次
Persistent=true
Unit=alpha-lab-a.service

[Install]
WantedBy=timers.target
