# research/deploy/systemd/alpha-lab-workflow.service
# alpha-lab workflow server 服務(長跑 background)
#
# 跑 `bun run workflow-server.ts`,在 background 持續 listen on
# WORKFLOW_SERVER_PORT(default 8090)。systemd timer (alpha-lab-a.timer,
# -b.timer, -d-pre.timer, -d-post.timer) 透過 curl 打到這個 service 內
# 對應的 /a /b /d/pre /d/post routes,以 trigger A/B/C/D workflow runs。
#
# 不依賴 systemd cron — 用 systemd timer 是 spec 決定(更可靠 + log 整合)

[Unit]
Description=alpha-lab Vercel Workflow runtime
After=network-online.target docker.service
Wants=network-online.target
# 需要 Postgres 跑起來(alpha-lab-postgres container 通常由 docker-compose up 帶起)
After=docker.service

[Service]
Type=simple
User=alpha
WorkingDirectory=/opt/alpha-lab/research
EnvironmentFile=/opt/alpha-lab/research/.env
ExecStart=/home/alpha/.local/share/mise/installs/bun/latest/bin/bun run workflow-server.ts
Restart=on-failure
RestartSec=10

# 讓 stdout 也寫 journal(就算 LOG_CONSOLE=false)
StandardOutput=journal
StandardError=journal

# 安全限制
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
