# research/deploy/systemd/alpha-lab-workflow.service
# alpha-lab workflow server 服務(長跑 background)
#
# 跑 `bun run workflow-server.ts`,在 background 持續 listen on
# WORKFLOW_SERVER_PORT(default 8090)。systemd timer (alpha-lab-a.timer,
# -b.timer, -d-pre.timer, -d-post.timer) 透過 curl 打到這個 service 內
# 對應的 /a /b /d/pre /d/post routes,以 trigger A/B/C/D workflow runs。
#
# 不依賴 systemd cron — 用 systemd timer 是 spec 決定(更可靠 + log 整合)
#
# 部署:__PUBLISH_USER__ placeholder 由 deploy-vm.sh sed 替換成實際 VM user
# (gcloud VM 預設 user 是 <account>_<domain>,e.g. q267009886_work_gmail_com)
#
# 安全限制:
#   ProtectSystem=strict 讓整個 FS read-only,但 log dir + drafts dir 需要寫入。
#   用 ReadWritePaths 白名單放行,比拿掉 ProtectSystem 更安全。
#   .well-known/ 是 workflow:build 產物,只在部署時寫(不是 runtime),不放行。
#   ProtectHome=read-only 讓 bun binary 在 ~/.bun/bin 可執行但不能改。

[Unit]
Description=alpha-lab Vercel Workflow runtime
After=network-online.target docker.service
Wants=network-online.target
# 需要 Postgres 跑起來(alpha-lab-postgres container 通常由 docker-compose up 帶起)
After=docker.service

[Service]
Type=simple
# User 跑時期由部署腳本 sed 改成實際 VM user(gcloud VM 預設非 alpha)
User=__PUBLISH_USER__
WorkingDirectory=/opt/alpha-lab/research
EnvironmentFile=/opt/alpha-lab/research/.env
# Bun 路徑:走 user 的 PATH(mise shims / ~/.bun/bin 都 cover)
Environment="PATH=/home/__PUBLISH_USER__/.bun/bin:/home/__PUBLISH_USER__/.local/share/mise/shims:/usr/local/bin:/usr/bin:/bin"
# BUN_TMPDIR / BUN_INSTALL:systemd 的預設 tmp namespace 跟 bun install cache
# 位置可能不一致(VM 實測 bun 1.3.x 報 'Unexpected accessing temporary directory'),
# 明確指定讓 bun 找到 cache + tmp 位置。即使拿掉 PrivateTmp=true 後,
# 攔截系統 /tmp 被 override 的情況(systemd service shell context 跟 user login shell 不同)
Environment="BUN_TMPDIR=/tmp"
Environment="BUN_INSTALL=/home/__PUBLISH_USER__/.bun"
ExecStart=bun run workflow-server.ts
Restart=on-failure
RestartSec=10

# 讓 stdout 也寫 journal(就算 LOG_CONSOLE=false)
StandardOutput=journal
StandardError=journal

# 安全限制 + 白名單
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# PrivateTmp=true 拿掉 — bun 1.3.x 在 private tmp namespace 內會找不到 global cache,
# 即使設 BUN_TMPDIR=/tmp 也只指到 private tmp 不是系統 / tmp,造成
# 'Unexpected accessing temporary directory. Please set $BUN_TMPDIR' 啟動失敗
# 對 alpha-lab workflow server 影響小(server 不處理 user input / file pipe)
# ReadWritePaths 白名單放行特定路徑,比拿掉 ProtectSystem 更安全
ReadWritePaths=/var/log/alpha-lab /opt/alpha-lab/research/logs /opt/alpha-lab/research/drafts /opt/alpha-lab/research/.tmp-bundles

[Install]
WantedBy=multi-user.target
