# automation/deploy/dagu/alpha-lab-dagu.service
[Unit]
Description=Alpha Lab unified Compose runtime
# 由 repository 維護的 canonical Compose 管理六個服務。
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/alpha-lab
ExecStart=/usr/bin/docker compose --env-file /etc/alpha-lab/stack.env -f /opt/alpha-lab/automation/deploy/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose --env-file /etc/alpha-lab/stack.env -f /opt/alpha-lab/automation/deploy/docker-compose.yml down
ExecReload=/usr/bin/docker compose --env-file /etc/alpha-lab/stack.env -f /opt/alpha-lab/automation/deploy/docker-compose.yml up -d --force-recreate --remove-orphans
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
