# automation/deploy/dagu/alpha-lab-dagu.service
[Unit]
Description=Alpha Lab Dagu runtime (compose wrapper)
# 此 unit 取代 legacy `dagu start-all` systemd unit，
# 驅動定義在 /etc/alpha-lab/docker-compose.yml 的
# Docker Compose stack。
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/etc/alpha-lab
ExecStart=/usr/bin/docker compose -f /etc/alpha-lab/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /etc/alpha-lab/docker-compose.yml down
ExecReload=/usr/bin/docker compose -f /etc/alpha-lab/docker-compose.yml up -d --force-recreate
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
